# webserver.service
[Unit]
Description=Web Server
After=network.target

[Service]
ExecStart=/usr/bin/webserver
Restart=on-failure
User=www-data
Group=www-data
WorkingDirectory=/var/www/html

[Install]
WantedBy=multi-user.target
